# deploy/autowrkers.service
[Unit]
Description=Autowrkers - Multi-session Claude Code Manager
After=network.target

[Service]
Type=simple
User=rat
Group=rat
WorkingDirectory=/home/rat/autowrkers
Environment="PATH=/home/rat/autowrkers/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/rat/autowrkers/venv/bin/python main.py start --host 0.0.0.0 --port 8420
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
